# files/system/usr/lib/systemd/user/jottad.service
[Unit]
Description=Jottacloud Daemon (via Distrobox)
After=network-online.target
Wants=network-online.target

[Service]
# This points to the binary Distrobox exported for you
ExecStart=%h/.local/bin/jottad
Restart=always
RestartSec=10
# This hides the daemon's window so it runs purely in the background
Environment="DISPLAY=" 

[Install]
WantedBy=default.target
